# Python/rocketmq/includes/app.conf
[program:app]
directory = /Users/a2htray/workspace/code-notebook/Python/rocketmq
command = python app.py
autostart = true
autorestart = true
user = a2htray
stdout_logfile = ./app.log